# backup_unit/backup.service
[Unit]
Description=Datensicherung mit rsync
[Service]
ExecStart=/usr/bin/rsync -a %h/Dokumente/Workdir/ /Backup/Dokumente/Workdir
